# roles/grafana/files/upstart.conf
description "grafana container"
author "Me"
start on filesystem and started nginx
stop on runlevel [!2345]
respawn
script
  /usr/bin/docker start -a grafana-container
end script