# roles/influxdb/files/upstart.conf
description "influxdb container"
author "Me"
start on filesystem and started nginx
stop on runlevel [!2345]
respawn
script
  /usr/bin/docker start -a influxdb-container
end script